# programs/gnome/autostart/dev.vencord.Vesktop.desktop
[Desktop Entry]
Type=Application
Name=dev.vencord.Vesktop
Exec=flatpak run dev.vencord.Vesktop
X-Flatpak=dev.vencord.Vesktop
