# programs/gnome/autostart/discord.desktop
[Desktop Entry]
Type=Application
Name=Discord
Exec=/usr/bin/discord
