# programs/gnome/autostart/discord.desktop
[Desktop Entry]
Type=Application
Name=Discord
Exec=flatpak run com.discordapp.Discord
X-Flatpak=com.discordapp.Discord
