# programs/gnome/autostart/firefox-nightly.desktop
[Desktop Entry]
Type=Application
Name=Firefox Nightly
Exec=/usr/bin/firefox-nightly
