# programs/gnome/autostart/firefox.desktop
[Desktop Entry]
Type=Application
Name=Firefox
Exec=/usr/bin/firefox
