# programs/gnome/autostart/gitcomet.desktop
[Desktop Entry]
Type=Application
Name=Gitcomet
Exec=/usr/bin/gitcomet
