# programs/gnome/autostart/intellij.desktop
[Desktop Entry]
Type=Application
Name=Intellij
Exec=/home/kahnwong/.local/share/JetBrains/Toolbox/apps/intellij-idea/bin/idea
