# programs/gnome/autostart/steam.desktop
[Desktop Entry]
Type=Application
Name=RustDesk
Exec=flatpak run com.valvesoftware.Steam
X-Flatpak=com.valvesoftware.Steam
