# programs/gnome/autostart/workspace-init.desktop
[Desktop Entry]
Type=Application
Name=Workspace Init
Exec=/home/kahnwong/go/bin/workspace-init
