# programs/autostart/ghostty.desktop
[Desktop Entry]
Type=Application
Name=Ghostty
Exec=/usr/bin/ghostty
