# programs/autostart/rustdesk.desktop
[Desktop Entry]
Type=Application
Name=RustDesk
Exec=flatpak run com.rustdesk.RustDesk
X-Flatpak=com.rustdesk.RustDesk
